# backup/.config/autostart/code.desktop
[Desktop Entry]
Exec=/usr/bin/code
Name=code
Type=Application
