# backup/.config/autostart/discord.desktop
[Desktop Entry]
Exec=/usr/bin/discord
Name=discord
Type=Application
