# backup/.config/autostart/firefox.desktop
[Desktop Entry]
Exec=/usr/bin/firefox
Name=firefox
Comment=a free and open source web browser from Mozilla
Type=Application
