# backup/.config/autostart/ibus-daemon.desktop
[Desktop Entry]
Exec=/usr/bin/ibus-daemon -drxR
Name=ibus-daemon
Comment=daemon program for ibus
Type=Application
